# monitoramento.service
[Unit]
Description=Monitoramento
[Service]
ExecStart=/usr/local/bin/node /home/patrick/app/service.js
